# deploy/supervisor.d/py-grammar-checker.ini
[program:py-grammar-checker]
process_name=py-grammar-checker
command=gunicorn py_grammar_checker:api -b :8200 -w 2
environment=PATH="/Users/rodmena/Code/english-learning-platform/src/micro-services/py-grammar-checker/env/bin:%(ENV_PATH)s"
user=rodmena
directory=/Users/rodmena/Code/english-learning-platform/src/micro-services/py-grammar-checker
group=staff
stdout_logfile =/Users/rodmena/Code/english-learning-platform/src/micro-services/py-grammar-checker/http.log   ; Where to write log messages
redirect_stderr = true
killasgroup=true
startsecs=5
stopwaitsecs=10
autostart=true

